# config/systemd/var-mnt-nvr.automount
[Unit]
Description= nvr automount

[Automount]
Where=/var/mnt/nvr

[Install]
WantedBy=multi-user.target
