# config/systemd/var-mnt-hdd-2TB\x2dWMC300210354.mount
[Unit]
Description=Mount local hdd

[Mount]
What=/dev/sda1
Where=/var/mnt/hdd/2TB-WMC300210354
Type=xfs

[Install]
WantedBy=default.target
