# config/systemd/var-mnt-hdd-2TB\x2dWMC300199477.mount
[Unit]
Description=Mount local hdd

[Mount]
What=/dev/sdb1
Where=/var/mnt/hdd/2TB-WMC300199477
Type=xfs

[Install]
WantedBy=default.target
